# systemd/.config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Unit]
After=graphical-session.target
PartOf=graphical-session.target
